# decommission-snapshot-20260426T003302Z/node-a/unit.service
[Unit]
Description=MXKeys cluster member node-a (raft smoke-test)
Documentation=file:///opt/matrix-family/test_servers/mxkeys-cluster/../README.md
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/matrix-family/test_servers/mxkeys-cluster/node-a
ExecStart=/opt/matrix-family/test_servers/mxkeys-cluster/bin/mxkeys -config /opt/matrix-family/test_servers/mxkeys-cluster/node-a/config.yaml
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=45
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mxkeys-node-a

[Install]
WantedBy=multi-user.target
